# arch/archbuild-catacomb.timer
[Unit]
Description=Minutely Catacomb package DB update

[Timer]
OnCalendar=*:*:0

[Install]
WantedBy=timers.target
