# arch/archbuild-catacomb.timer
[Unit]
Description=5-Minutely Catacomb package DB update

[Timer]
OnCalendar=*:*:0/5

[Install]
WantedBy=timers.target
