# arch/archbuild-catacomb.timer
[Unit]
Description=15-Minutely Catacomb package DB update

[Timer]
OnCalendar=*:0/15:00

[Install]
WantedBy=timers.target
