# arch/hexagonrpcd/hexagonrpcd-sdsp.service
[Unit]
Description=Hexagonrpcd SDSP Daemon to support Qualcomm Hexagon SDSP virtual filesystem
Requires=dev-fastrpc\x2dsdsp.device
After=dev-fastrpc\x2dsdsp.device

[Service]
User=fastrpc
Group=fastrpc
EnvironmentFile=/etc/conf.d/hexagonrpcd-sdsp
ExecStart=/usr/bin/hexagonrpcd -f /dev/fastrpc-sdsp -d sdsp -s -R ${hexagonrpcd_fw_dir}

[Install]
WantedBy=multi-user.target
