# arch/hexagonrpcd/hexagonrpcd-adsp-sensorspd.service
[Unit]
Description=Hexagonrpcd ADSP SensorPD Daemon to support Qualcomm Hexagon ADSP virtual filesystem for SensorPD
Requires=dev-fastrpc\x2dadsp.device
After=dev-fastrpc\x2dadsp.device

[Service]
User=fastrpc
Group=fastrpc
EnvironmentFile=/etc/conf.d/hexagonrpcd-adsp-sensorpd
ExecStart=/usr/bin/hexagonrpcd -f /dev/fastrpc-adsp -d adsp -R ${hexagonrpcd_fw_dir}

[Install]
WantedBy=multi-user.target
